# overlay/etc/profile.d/emcomm-tools.sh
export PATH=/opt/emcomm-tools/bin:$PATH
